# var/lib/dpkg/info/erlang-base.list
/.
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/epmd.service
/lib/systemd/system/epmd.socket
/usr
/usr/bin
/usr/lib
/usr/lib/erlang
/usr/lib/erlang/bin
/usr/lib/erlang/bin/erl
/usr/lib/erlang/bin/erlc
/usr/lib/erlang/bin/escript
/usr/lib/erlang/bin/no_dot_erlang.boot
/usr/lib/erlang/bin/run_erl
/usr/lib/erlang/bin/start
/usr/lib/erlang/bin/start.boot
/usr/lib/erlang/bin/start.script
/usr/lib/erlang/bin/start_clean.boot
/usr/lib/erlang/bin/start_erl
/usr/lib/erlang/bin/start_sasl.boot
/usr/lib/erlang/bin/to_erl
/usr/lib/erlang/erts-8.2.1
/usr/lib/erlang/erts-8.2.1/bin
/usr/lib/erlang/erts-8.2.1/bin/beam
/usr/lib/erlang/erts-8.2.1/bin/beam.smp
/usr/lib/erlang/erts-8.2.1/bin/dyn_erl
/usr/lib/erlang/erts-8.2.1/bin/epmd
/usr/lib/erlang/erts-8.2.1/bin/erl
/usr/lib/erlang/erts-8.2.1/bin/erl_child_setup
/usr/lib/erlang/erts-8.2.1/bin/erlc
/usr/lib/erlang/erts-8.2.1/bin/erlexec
/usr/lib/erlang/erts-8.2.1/bin/escript
/usr/lib/erlang/erts-8.2.1/bin/heart
/usr/lib/erlang/erts-8.2.1/bin/inet_gethost
/usr/lib/erlang/erts-8.2.1/bin/run_erl
/usr/lib/erlang/erts-8.2.1/bin/start
/usr/lib/erlang/erts-8.2.1/bin/start_erl
/usr/lib/erlang/erts-8.2.1/bin/to_erl
/usr/lib/erlang/lib
/usr/lib/erlang/lib/compiler-7.0.3
/usr/lib/erlang/lib/compiler-7.0.3/ebin
/usr/lib/erlang/lib/compiler-7.0.3/ebin/beam_a.beam
/usr/lib/erlang/lib/compiler-7.0.3/ebin/beam_asm.beam
/usr/lib/erlang/lib/compiler-7.0.3/ebin/beam_block.beam
/usr/lib/erlang/lib/compiler-7.0.3/ebin/beam_bool.beam
/usr/lib/erlang/lib/compiler-7.0.3/ebin/beam_bs.beam
/usr/lib/erlang/lib/compiler-7.0.3/ebin/beam_bsm.beam
/usr/lib/erlang/lib/compiler-7.0.3/ebin/beam_clean.beam
/usr/lib/erlang/lib/compiler-7.0.3/ebin/beam_dead.beam
/usr/lib/erlang/lib/compiler-7.0.3/ebin/beam_dict.beam
/usr/lib/erlang/lib/compiler-7.0.3/ebin/beam_disasm.beam
/usr/lib/erlang/lib/compiler-7.0.3/ebin/beam_except.beam
/usr/lib/erlang/lib/compiler-7.0.3/ebin/beam_flatten.beam
/usr/lib/erlang/lib/compiler-7.0.3/ebin/beam_jump.beam
/usr/lib/erlang/lib/compiler-7.0.3/ebin/beam_listing.beam
/usr/lib/erlang/lib/compiler-7.0.3/ebin/beam_opcodes.beam
/usr/lib/erlang/lib/compiler-7.0.3/ebin/beam_peep.beam
/usr/lib/erlang/lib/compiler-7.0.3/ebin/beam_receive.beam
/usr/lib/erlang/lib/compiler-7.0.3/ebin/beam_reorder.beam
/usr/lib/erlang/lib/compiler-7.0.3/ebin/beam_split.beam
/usr/lib/erlang/lib/compiler-7.0.3/ebin/beam_trim.beam
/usr/lib/erlang/lib/compiler-7.0.3/ebin/beam_type.beam
/usr/lib/erlang/lib/compiler-7.0.3/ebin/beam_utils.beam
/usr/lib/erlang/lib/compiler-7.0.3/ebin/beam_validator.beam
/usr/lib/erlang/lib/compiler-7.0.3/ebin/beam_z.beam
/usr/lib/erlang/lib/compiler-7.0.3/ebin/cerl.beam
/usr/lib/erlang/lib/compiler-7.0.3/ebin/cerl_clauses.beam
/usr/lib/erlang/lib/compiler-7.0.3/ebin/cerl_inline.beam
/usr/lib/erlang/lib/compiler-7.0.3/ebin/cerl_sets.beam
/usr/lib/erlang/lib/compiler-7.0.3/ebin/cerl_trees.beam
/usr/lib/erlang/lib/compiler-7.0.3/ebin/compile.beam
/usr/lib/erlang/lib/compiler-7.0.3/ebin/compiler.app
/usr/lib/erlang/lib/compiler-7.0.3/ebin/compiler.appup
/usr/lib/erlang/lib/compiler-7.0.3/ebin/core_lib.beam
/usr/lib/erlang/lib/compiler-7.0.3/ebin/core_lint.beam
/usr/lib/erlang/lib/compiler-7.0.3/ebin/core_parse.beam
/usr/lib/erlang/lib/compiler-7.0.3/ebin/core_pp.beam
/usr/lib/erlang/lib/compiler-7.0.3/ebin/core_scan.beam
/usr/lib/erlang/lib/compiler-7.0.3/ebin/erl_bifs.beam
/usr/lib/erlang/lib/compiler-7.0.3/ebin/rec_env.beam
/usr/lib/erlang/lib/compiler-7.0.3/ebin/sys_core_dsetel.beam
/usr/lib/erlang/lib/compiler-7.0.3/ebin/sys_core_fold.beam
/usr/lib/erlang/lib/compiler-7.0.3/ebin/sys_core_fold_lists.beam
/usr/lib/erlang/lib/compiler-7.0.3/ebin/sys_core_inline.beam
/usr/lib/erlang/lib/compiler-7.0.3/ebin/sys_pre_attributes.beam
/usr/lib/erlang/lib/compiler-7.0.3/ebin/sys_pre_expand.beam
/usr/lib/erlang/lib/compiler-7.0.3/ebin/v3_codegen.beam
/usr/lib/erlang/lib/compiler-7.0.3/ebin/v3_core.beam
/usr/lib/erlang/lib/compiler-7.0.3/ebin/v3_kernel.beam
/usr/lib/erlang/lib/compiler-7.0.3/ebin/v3_kernel_pp.beam
/usr/lib/erlang/lib/compiler-7.0.3/ebin/v3_life.beam
/usr/lib/erlang/lib/erl_interface-3.9.2
/usr/lib/erlang/lib/erl_interface-3.9.2/bin
/usr/lib/erlang/lib/erl_interface-3.9.2/bin/erl_call
/usr/lib/erlang/lib/erts-8.2.1
/usr/lib/erlang/lib/erts-8.2.1/ebin
/usr/lib/erlang/lib/erts-8.2.1/ebin/erl_prim_loader.beam
/usr/lib/erlang/lib/erts-8.2.1/ebin/erl_tracer.beam
/usr/lib/erlang/lib/erts-8.2.1/ebin/erlang.beam
/usr/lib/erlang/lib/erts-8.2.1/ebin/erts.app
/usr/lib/erlang/lib/erts-8.2.1/ebin/erts_code_purger.beam
/usr/lib/erlang/lib/erts-8.2.1/ebin/erts_dirty_process_code_checker.beam
/usr/lib/erlang/lib/erts-8.2.1/ebin/erts_internal.beam
/usr/lib/erlang/lib/erts-8.2.1/ebin/erts_literal_area_collector.beam
/usr/lib/erlang/lib/erts-8.2.1/ebin/init.beam
/usr/lib/erlang/lib/erts-8.2.1/ebin/otp_ring0.beam
/usr/lib/erlang/lib/erts-8.2.1/ebin/prim_eval.beam
/usr/lib/erlang/lib/erts-8.2.1/ebin/prim_file.beam
/usr/lib/erlang/lib/erts-8.2.1/ebin/prim_inet.beam
/usr/lib/erlang/lib/erts-8.2.1/ebin/prim_zip.beam
/usr/lib/erlang/lib/erts-8.2.1/ebin/zlib.beam
/usr/lib/erlang/lib/hipe-3.15.3
/usr/lib/erlang/lib/hipe-3.15.3/ebin
/usr/lib/erlang/lib/hipe-3.15.3/ebin/cerl_cconv.beam
/usr/lib/erlang/lib/hipe-3.15.3/ebin/cerl_closurean.beam
/usr/lib/erlang/lib/hipe-3.15.3/ebin/cerl_hipeify.beam
/usr/lib/erlang/lib/hipe-3.15.3/ebin/cerl_lib.beam
/usr/lib/erlang/lib/hipe-3.15.3/ebin/cerl_messagean.beam
/usr/lib/erlang/lib/hipe-3.15.3/ebin/cerl_pmatch.beam
/usr/lib/erlang/lib/hipe-3.15.3/ebin/cerl_prettypr.beam
/usr/lib/erlang/lib/hipe-3.15.3/ebin/cerl_to_icode.beam
/usr/lib/erlang/lib/hipe-3.15.3/ebin/cerl_typean.beam
/usr/lib/erlang/lib/hipe-3.15.3/ebin/erl_bif_types.beam
/usr/lib/erlang/lib/hipe-3.15.3/ebin/erl_types.beam
/usr/lib/erlang/lib/hipe-3.15.3/ebin/hipe.app
/usr/lib/erlang/lib/hipe-3.15.3/ebin/hipe.appup
/usr/lib/erlang/lib/hipe-3.15.3/ebin/hipe_bb.beam
/usr/lib/erlang/lib/hipe-3.15.3/ebin/hipe_beam_to_icode.beam
/usr/lib/erlang/lib/hipe-3.15.3/ebin/hipe_consttab.beam
/usr/lib/erlang/lib/hipe-3.15.3/ebin/hipe_digraph.beam
/usr/lib/erlang/lib/hipe-3.15.3/ebin/hipe_dominators.beam
/usr/lib/erlang/lib/hipe-3.15.3/ebin/hipe_dot.beam
/usr/lib/erlang/lib/hipe-3.15.3/ebin/hipe_gen_cfg.beam
/usr/lib/erlang/lib/hipe-3.15.3/ebin/hipe_gensym.beam
/usr/lib/erlang/lib/hipe-3.15.3/ebin/hipe_icode.beam
/usr/lib/erlang/lib/hipe-3.15.3/ebin/hipe_icode_bincomp.beam
/usr/lib/erlang/lib/hipe-3.15.3/ebin/hipe_icode_call_elim.beam
/usr/lib/erlang/lib/hipe-3.15.3/ebin/hipe_icode_callgraph.beam
/usr/lib/erlang/lib/hipe-3.15.3/ebin/hipe_icode_cfg.beam
/usr/lib/erlang/lib/hipe-3.15.3/ebin/hipe_icode_coordinator.beam
/usr/lib/erlang/lib/hipe-3.15.3/ebin/hipe_icode_ebb.beam
/usr/lib/erlang/lib/hipe-3.15.3/ebin/hipe_icode_exceptions.beam
/usr/lib/erlang/lib/hipe-3.15.3/ebin/hipe_icode_fp.beam
/usr/lib/erlang/lib/hipe-3.15.3/ebin/hipe_icode_inline_bifs.beam
/usr/lib/erlang/lib/hipe-3.15.3/ebin/hipe_icode_instruction_counter.beam
/usr/lib/erlang/lib/hipe-3.15.3/ebin/hipe_icode_liveness.beam
/usr/lib/erlang/lib/hipe-3.15.3/ebin/hipe_icode_mulret.beam
/usr/lib/erlang/lib/hipe-3.15.3/ebin/hipe_icode_pp.beam
/usr/lib/erlang/lib/hipe-3.15.3/ebin/hipe_icode_primops.beam
/usr/lib/erlang/lib/hipe-3.15.3/ebin/hipe_icode_range.beam
/usr/lib/erlang/lib/hipe-3.15.3/ebin/hipe_icode_split_arith.beam
/usr/lib/erlang/lib/hipe-3.15.3/ebin/hipe_icode_ssa.beam
/usr/lib/erlang/lib/hipe-3.15.3/ebin/hipe_icode_ssa_const_prop.beam
/usr/lib/erlang/lib/hipe-3.15.3/ebin/hipe_icode_ssa_copy_prop.beam
/usr/lib/erlang/lib/hipe-3.15.3/ebin/hipe_icode_ssa_struct_reuse.beam
/usr/lib/erlang/lib/hipe-3.15.3/ebin/hipe_icode_type.beam
/usr/lib/erlang/lib/hipe-3.15.3/ebin/hipe_main.beam
/usr/lib/erlang/lib/hipe-3.15.3/ebin/hipe_timing.beam
/usr/lib/erlang/lib/kernel-5.1.1
/usr/lib/erlang/lib/kernel-5.1.1/ebin
/usr/lib/erlang/lib/kernel-5.1.1/ebin/application.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/application_controller.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/application_master.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/application_starter.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/auth.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/code.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/code_server.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/disk_log.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/disk_log_1.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/disk_log_server.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/disk_log_sup.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/dist_ac.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/dist_util.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/erl_boot_server.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/erl_ddll.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/erl_distribution.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/erl_epmd.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/erl_reply.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/error_handler.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/error_logger.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/erts_debug.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/file.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/file_io_server.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/file_server.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/gen_sctp.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/gen_tcp.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/gen_udp.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/global.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/global_group.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/global_search.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/group.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/heart.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/hipe_unified_loader.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/inet.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/inet6_sctp.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/inet6_tcp.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/inet6_tcp_dist.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/inet6_udp.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/inet_config.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/inet_db.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/inet_dns.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/inet_gethost_native.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/inet_hosts.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/inet_parse.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/inet_res.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/inet_sctp.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/inet_tcp.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/inet_tcp_dist.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/inet_udp.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/kernel.app
/usr/lib/erlang/lib/kernel-5.1.1/ebin/kernel.appup
/usr/lib/erlang/lib/kernel-5.1.1/ebin/kernel.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/kernel_config.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/local_tcp.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/local_udp.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/net.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/net_adm.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/net_kernel.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/os.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/pg2.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/ram_file.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/rpc.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/seq_trace.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/standard_error.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/user.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/user_drv.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/user_sup.beam
/usr/lib/erlang/lib/kernel-5.1.1/ebin/wrap_log_reader.beam
/usr/lib/erlang/lib/kernel-5.1.1/include
/usr/lib/erlang/lib/kernel-5.1.1/include/dist.hrl
/usr/lib/erlang/lib/kernel-5.1.1/include/dist_util.hrl
/usr/lib/erlang/lib/kernel-5.1.1/include/file.hrl
/usr/lib/erlang/lib/kernel-5.1.1/include/inet.hrl
/usr/lib/erlang/lib/kernel-5.1.1/include/inet_sctp.hrl
/usr/lib/erlang/lib/kernel-5.1.1/include/net_address.hrl
/usr/lib/erlang/lib/sasl-3.0.2
/usr/lib/erlang/lib/sasl-3.0.2/ebin
/usr/lib/erlang/lib/sasl-3.0.2/ebin/alarm_handler.beam
/usr/lib/erlang/lib/sasl-3.0.2/ebin/erlsrv.beam
/usr/lib/erlang/lib/sasl-3.0.2/ebin/format_lib_supp.beam
/usr/lib/erlang/lib/sasl-3.0.2/ebin/misc_supp.beam
/usr/lib/erlang/lib/sasl-3.0.2/ebin/rb.beam
/usr/lib/erlang/lib/sasl-3.0.2/ebin/rb_format_supp.beam
/usr/lib/erlang/lib/sasl-3.0.2/ebin/release_handler.beam
/usr/lib/erlang/lib/sasl-3.0.2/ebin/release_handler_1.beam
/usr/lib/erlang/lib/sasl-3.0.2/ebin/sasl.app
/usr/lib/erlang/lib/sasl-3.0.2/ebin/sasl.appup
/usr/lib/erlang/lib/sasl-3.0.2/ebin/sasl.beam
/usr/lib/erlang/lib/sasl-3.0.2/ebin/sasl_report.beam
/usr/lib/erlang/lib/sasl-3.0.2/ebin/sasl_report_file_h.beam
/usr/lib/erlang/lib/sasl-3.0.2/ebin/sasl_report_tty_h.beam
/usr/lib/erlang/lib/sasl-3.0.2/ebin/si.beam
/usr/lib/erlang/lib/sasl-3.0.2/ebin/si_sasl_supp.beam
/usr/lib/erlang/lib/sasl-3.0.2/ebin/systools.beam
/usr/lib/erlang/lib/sasl-3.0.2/ebin/systools_lib.beam
/usr/lib/erlang/lib/sasl-3.0.2/ebin/systools_make.beam
/usr/lib/erlang/lib/sasl-3.0.2/ebin/systools_rc.beam
/usr/lib/erlang/lib/sasl-3.0.2/ebin/systools_relup.beam
/usr/lib/erlang/lib/sasl-3.0.2/examples
/usr/lib/erlang/lib/sasl-3.0.2/examples/ebin
/usr/lib/erlang/lib/sasl-3.0.2/examples/ebin/target_system.beam
/usr/lib/erlang/lib/stdlib-3.2
/usr/lib/erlang/lib/stdlib-3.2/ebin
/usr/lib/erlang/lib/stdlib-3.2/ebin/array.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/base64.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/beam_lib.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/binary.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/c.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/calendar.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/dets.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/dets_server.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/dets_sup.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/dets_utils.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/dets_v8.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/dets_v9.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/dict.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/digraph.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/digraph_utils.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/edlin.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/edlin_expand.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/epp.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/erl_anno.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/erl_bits.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/erl_compile.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/erl_eval.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/erl_expand_records.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/erl_internal.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/erl_lint.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/erl_parse.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/erl_posix_msg.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/erl_pp.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/erl_scan.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/erl_tar.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/error_logger_file_h.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/error_logger_tty_h.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/escript.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/ets.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/eval_bits.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/file_sorter.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/filelib.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/filename.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/gb_sets.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/gb_trees.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/gen.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/gen_event.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/gen_fsm.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/gen_server.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/gen_statem.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/io.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/io_lib.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/io_lib_format.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/io_lib_fread.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/io_lib_pretty.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/lib.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/lists.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/log_mf_h.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/maps.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/math.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/ms_transform.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/orddict.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/ordsets.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/otp_internal.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/pool.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/proc_lib.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/proplists.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/qlc.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/qlc_pt.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/queue.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/rand.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/random.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/re.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/sets.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/shell.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/shell_default.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/slave.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/sofs.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/stdlib.app
/usr/lib/erlang/lib/stdlib-3.2/ebin/stdlib.appup
/usr/lib/erlang/lib/stdlib-3.2/ebin/string.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/supervisor.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/supervisor_bridge.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/sys.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/timer.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/unicode.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/win32reg.beam
/usr/lib/erlang/lib/stdlib-3.2/ebin/zip.beam
/usr/lib/erlang/lib/stdlib-3.2/include
/usr/lib/erlang/lib/stdlib-3.2/include/assert.hrl
/usr/lib/erlang/lib/stdlib-3.2/include/erl_bits.hrl
/usr/lib/erlang/lib/stdlib-3.2/include/erl_compile.hrl
/usr/lib/erlang/lib/stdlib-3.2/include/ms_transform.hrl
/usr/lib/erlang/lib/stdlib-3.2/include/qlc.hrl
/usr/lib/erlang/lib/stdlib-3.2/include/zip.hrl
/usr/lib/erlang/releases
/usr/lib/erlang/releases/19
/usr/lib/erlang/releases/19/OTP_VERSION
/usr/lib/erlang/releases/19/installed_application_versions
/usr/lib/erlang/releases/19/no_dot_erlang.boot
/usr/lib/erlang/releases/19/no_dot_erlang.rel
/usr/lib/erlang/releases/19/no_dot_erlang.script
/usr/lib/erlang/releases/19/start.boot
/usr/lib/erlang/releases/19/start.script
/usr/lib/erlang/releases/19/start_all_example.rel
/usr/lib/erlang/releases/19/start_clean.boot
/usr/lib/erlang/releases/19/start_clean.rel
/usr/lib/erlang/releases/19/start_clean.script
/usr/lib/erlang/releases/19/start_sasl.boot
/usr/lib/erlang/releases/19/start_sasl.rel
/usr/lib/erlang/releases/19/start_sasl.script
/usr/lib/erlang/releases/RELEASES
/usr/lib/erlang/releases/start_erl.data
/usr/share
/usr/share/apps
/usr/share/apps/konsole
/usr/share/apps/konsole/erlang.desktop
/usr/share/doc
/usr/share/doc/erlang-base
/usr/share/doc/erlang-base/AUTHORS.gz
/usr/share/doc/erlang-base/README.Debian.gz
/usr/share/doc/erlang-base/TODO.Debian
/usr/share/doc/erlang-base/changelog.Debian.gz
/usr/share/doc/erlang-base/copyright
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/erlang-base
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/epmd.1.gz
/usr/share/man/man1/erl.1.gz
/usr/share/man/man1/erl_call.1.gz
/usr/share/man/man1/erlc.1.gz
/usr/share/man/man1/escript.1.gz
/usr/share/man/man1/run_erl.1.gz
/usr/share/man/man1/start_embedded.1.gz
/usr/share/man/man1/to_erl.1.gz
/usr/share/menu
/usr/share/menu/erlang-base
/usr/share/pixmaps
/usr/share/pixmaps/erlang.xpm
/usr/bin/epmd
/usr/bin/erl
/usr/bin/erl_call
/usr/bin/erlc
/usr/bin/escript
/usr/bin/run_erl
/usr/bin/start_embedded
/usr/bin/to_erl
/usr/lib/erlang/bin/epmd
